# files/system/etc/profile.d/default-editor.sh
# DEFAULT EDITOR SHOULD BE VIM, NOT NANO

if [ -z "$EDITOR" ]; then
    export EDITOR=vim
fi